# mmfw-sysconf/mmfw-sysconf-sc9830/etc/profile.d/mmf.sh
export GST_DEBUG=2
export PULSE_CLIENTCONFIG=/etc/pulse/client.conf
export SPRD_DEBUG=2
